# config/supervisord.conf
[unix_http_server]
file=/tmp/claw-supervisor.sock

[supervisord]
logfile=logs/supervisord.log
logfile_maxbytes=10MB
logfile_backups=3
loglevel=info
pidfile=/tmp/claw-supervisord.pid
nodaemon=false

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/claw-supervisor.sock

[program:runner]
command=venv/bin/python -m app.runner
directory=/Users/henry_oc/develop/claw-trader
environment=PYTHONPATH="src",PYTHONUNBUFFERED="1"
stdout_logfile=logs/runner.log
stderr_logfile=logs/runner.log
stdout_logfile_maxbytes=10MB
autostart=true
autorestart=true
startretries=5
stopwaitsecs=10

[program:market_data_runner]
command=venv/bin/python -m app.market_data_runner
directory=/Users/henry_oc/develop/claw-trader
environment=PYTHONPATH="src",PYTHONUNBUFFERED="1"
stdout_logfile=logs/market_data.log
stderr_logfile=logs/market_data.log
stdout_logfile_maxbytes=10MB
autostart=true
autorestart=true
startretries=5
stopwaitsecs=10

[program:order_watcher]
command=venv/bin/python -m scripts.order_watcher
directory=/Users/henry_oc/develop/claw-trader
environment=PYTHONPATH="src",PYTHONUNBUFFERED="1",WATCHER_TTL_CANCEL_SEC="60"
stdout_logfile=logs/order_watcher.log
stderr_logfile=logs/order_watcher.log
stdout_logfile_maxbytes=10MB
autostart=true
autorestart=true
startretries=5
stopwaitsecs=10

[program:ai_dual_eval_runner]
command=venv/bin/python -m app.ai_dual_eval_runner
directory=/Users/henry_oc/develop/claw-trader
environment=PYTHONPATH="src",PYTHONUNBUFFERED="1"
stdout_logfile=logs/ai_dual_eval.log
stderr_logfile=logs/ai_dual_eval.log
stdout_logfile_maxbytes=10MB
autostart=true
autorestart=true
startretries=5
stopwaitsecs=10

[program:consensus_signal_runner]
command=venv/bin/python -m app.consensus_signal_runner
directory=/Users/henry_oc/develop/claw-trader
environment=PYTHONPATH="src",PYTHONUNBUFFERED="1"
stdout_logfile=logs/consensus_signal.log
stderr_logfile=logs/consensus_signal.log
stdout_logfile_maxbytes=10MB
autostart=true
autorestart=true
startretries=5
stopwaitsecs=10

[program:openclaw_bot]
command=venv/bin/python -m app.openclaw_bot
directory=/Users/henry_oc/develop/claw-trader
environment=PYTHONPATH="src",PYTHONUNBUFFERED="1"
stdout_logfile=logs/openclaw_bot.log
stderr_logfile=logs/openclaw_bot.log
stdout_logfile_maxbytes=10MB
autostart=true
autorestart=true
startretries=5
stopwaitsecs=10

[program:news_runner]
command=venv/bin/python -m app.news_runner
directory=/Users/henry_oc/develop/claw-trader
environment=PYTHONPATH="src",PYTHONUNBUFFERED="1"
stdout_logfile=logs/news_runner.log
stderr_logfile=logs/news_runner.log
stdout_logfile_maxbytes=10MB
autostart=true
autorestart=true
startretries=5
stopwaitsecs=10

[program:position_exit_runner]
command=venv/bin/python -m scripts.position_exit_runner
directory=/Users/henry_oc/develop/claw-trader
environment=PYTHONPATH="src",PYTHONUNBUFFERED="1"
stdout_logfile=logs/position_exit.log
stderr_logfile=logs/position_exit.log
stdout_logfile_maxbytes=10MB
autostart=true
autorestart=true
startretries=5
stopwaitsecs=10

[program:position_engine]
command=venv/bin/python -m scripts.position_engine
directory=/Users/henry_oc/develop/claw-trader
environment=PYTHONPATH="src",PYTHONUNBUFFERED="1"
stdout_logfile=logs/position_engine.log
stderr_logfile=logs/position_engine.log
stdout_logfile_maxbytes=10MB
autostart=true
autorestart=true
startretries=5
stopwaitsecs=10

[program:watchlist_selector]
command=venv/bin/python -m scripts.watchlist_selector_runner
directory=/Users/henry_oc/develop/claw-trader
environment=PYTHONPATH="src",PYTHONUNBUFFERED="1"
stdout_logfile=logs/watchlist_selector.log
stderr_logfile=logs/watchlist_selector.log
stdout_logfile_maxbytes=10MB
autostart=true
autorestart=true
startretries=5
stopwaitsecs=10

[program:daily_report_runner]
command=venv/bin/python -m scripts.daily_report_runner
directory=/Users/henry_oc/develop/claw-trader
environment=PYTHONPATH="src",PYTHONUNBUFFERED="1"
stdout_logfile=logs/daily_report.log
stderr_logfile=logs/daily_report.log
stdout_logfile_maxbytes=10MB
autostart=true
autorestart=true
startretries=5
stopwaitsecs=10

[eventlistener:crash-notifier]
command=venv/bin/python -m scripts.supervisor_crash_notifier
events=PROCESS_STATE_FATAL
directory=/Users/henry_oc/develop/claw-trader
environment=PYTHONPATH="src"
